# documentos/service Restful (init.d)
#!/bin/sh -e
#/etc/init.d/Restful

### BEGIN INIT INFO
# Provides:                     Restful
# Required-Start:               $local_fs $remote_fs $network
# Required-Stop:                $local_fs $remote_fs $network
# Default-Start:                2 3 4 5
# Default-Stop:                 0 1 6
# Short-Description:            Inicia api restful
# Description:                  Inicia Servidor RESTful.
### END INIT INFO

# Script de inicializacao LSB (https://wiki.debian.org/LSBInitScripts)
# Autor: Erivando Sena
# Copyright: (c) 2015 Unilab
# Email: erivandoramos@unilab.edu.br
# Data: 22/07/2015
# Status: Prototype  (Prototype | Development | Production)

# Para INSTALAR:
# (DESCONTINUADO partir do Debian 6.0) sudo update-rc.d Restful defaults
# (DESCONTINUADO partir do Debian 6.0) sudo update-rc.d -f Restful remove
# sudo insserv -v Restful
# sudo insserv -v Restful -r
# Verificar na lista de servicos ativos: sudo service --status-all

DIR=/home/erivando
DAEMON=$DIR/run_restful.sh
DAEMON_NAME=Restful
DAEMON_USER=root
PIDFILE=/var/run/$DAEMON_NAME.pid
DAEMON_OPTS=""

. /lib/lsb/init-functions

do_start() {
  log_daemon_msg "Iniciando $DAEMON_NAME"
  start-stop-daemon --verbose --start --pidfile $PIDFILE --make-pidfile --user $DAEMON_USER --chuid $DAEMON_USER --startas $DAEMON -- $DAEMON_OPTS
  log_end_msg $?                     #--background
}
do_stop() {
  log_daemon_msg "Parando $DAEMON_NAME"
  start-stop-daemon --stop --pidfile $PIDFILE --retry 10
  log_end_msg $?
}
case "$1" in

  start|stop)
    do_${1}
    ;;

  restart|reload|force-reload)
    do_stop
    do_start
    ;;

  status)
    status_of_proc "$DAEMON_NAME" "$DAEMON" && exit 0 || exit $?
    ;;

  *)
    echo "Use: /etc/init.d/$DEAMON_NAME {start|stop|restart|status}"
    echo "Ou use:"
    echo "service ServidorRestful start"
    echo "service ServidorRestful stop"
    echo "service ServidorRestful restart"
    echo "service ServidorRestful status"
    exit 1
    ;;

esac
exit 0
